# Backup/db-backup.timer
[Unit]
Description=Timer to trigger the RDS DB Backup Service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
