# Backup/db-backup.service
[Unit]
Description=RDS Database Backup Service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dbbackup
User=root
Group=root
