# backup/dbbackup.service
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/dbbackup serpent_surge_db score
